# scripts/spectra.service
[Unit]
Description=SPECTRA Core API
After=network.target

[Service]
Type=simple
User=u
WorkingDirectory=/home/u/dev/spectra
Environment="PATH=/home/u/dev/spectra/.venv/bin"
EnvironmentFile=/home/u/dev/spectra/.env
ExecStart=/home/u/dev/spectra/.venv/bin/uvicorn core.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
