# scripts/spectra-tunnel.service
[Unit]
Description=SPECTRA Cloudflare Tunnel
After=network.target spectra.service

[Service]
Type=simple
User=u
ExecStart=/usr/bin/cloudflared tunnel run spectra
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
